# stats.service
[Unit]
Description=SSD1306 Display

[Service]
WorkingDirectory=/opt/stats
ExecStart=/opt/stats/stats.py

[Install]
WantedBy=multi-user.target
